# config/multirotor/quad/ppd/target/airliner.service
[Unit]
Description=Airliner UAS Flight Software
After=network.target auditd.service
ConditionPathExists=!/opt/airliner/disable

[Service]
Type=simple
WorkingDirectory=/opt/airliner/
ExecStart=/opt/airliner/airliner
ExecReload=/bin/kill -9 $MAINPID
#ExecReload=/opt/airliner/allstop
ExecStop=/bin/kill -9 $MAINPID
KillMode=process
#Restart=on-failure

[Install]
WantedBy=multi-user.target
Alias=airliner.service
